# packaging/systemd/codex-mobile-safe.service.in
[Unit]
Description=Codex Mobile Safe (Tailnet-only backend)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@PROJECT_DIR@
ExecStart=@PREFIX@/bin/codex-mobile-safe start @PROJECT_DIR@ --port 5900 --password-file %h/.codex/codex-mobile-safe-password --no-open --no-login --sandbox-mode danger-full-access --approval-policy never
Restart=on-failure
RestartSec=3
UMask=0077

[Install]
WantedBy=default.target
